# systemd/ubus.socket.in
[Unit]
Description=OpenWrt micro bus socket

[Socket]
ListenStream=@UBUS_UNIX_SOCKET@

[Install]
WantedBy=sockets.target
