# systemd/ubus.service.in
[Unit]
Description=OpenWrt micro bus
Requires=ubus.socket

[Service]
ExecStart=@UBUSD_BINARY@
